# etc/supervisor/conf.d/go.conf
[program:go]
command=bash -c '/dockerstartup/go_go.sh'
#/dockerstartup/startup.sh
#sudo  vncserver -localhost  -depth 24 -geometry 1360x768 :1
autorestart=false
directory=/root